# backup/restic-prune.timer
[Unit]
Description=Wöchentliches Restic Prune

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
